# debfiles/sinagios.postinstall
# Have to do this since FPM doesn't know how to set ownership on files
/bin/chown -R sinagios.sinagios /var/log/sinagios && update-rc.d sinagios defaults
